# systemd/network-detection-check.service
[Unit]
Description=Network detection checker

[Service]
ExecStart=/etc/monitor-scripts/network-detection.sh
